# deploy/systemd/lgm-receiver.service
[Unit]
Description=LGM Receiver Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/lgm-receiver --config /etc/lgm-monitor/config.json
Restart=always
RestartSec=5
User=root
Group=root
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/etc/lgm-monitor /var/lib/lgm-monitor /var/log/lgm-monitor

[Install]
WantedBy=multi-user.target
